# etc/supervisor/conf.d/beacon.conf
[program:beacon]
command = python /var/www/router/bin/broadcast_client.py join
autostart = true
startsecs = 5
user = pi
redirect_stderr = true
